# rpa-websocket.service
# Systemd service para WebSocket Server
# Salvar como: /etc/systemd/system/rpa-websocket.service

[Unit]
Description=RPA WebSocket Server
After=network.target redis.service
Requires=redis.service

[Service]
Type=simple
User=www-data
Group=www-data
WorkingDirectory=/var/www/rpaimediatoseguros.com.br/websocket
ExecStart=/usr/bin/node server.js
Restart=always
RestartSec=5
StandardOutput=syslog
StandardError=syslog
SyslogIdentifier=rpa-websocket

# Variáveis de ambiente
Environment=NODE_ENV=production
Environment=WS_PORT=8080
Environment=REDIS_HOST=localhost
Environment=REDIS_PORT=6379

[Install]
WantedBy=multi-user.target








